# notify-listener.service
[Unit]
Description=irssi-libnotify-git

[Service]
ExecStart=/usr/local/bin/notify-listener.py

[Install]
WantedBy=default.target
